# scripts/crontab
# Repository fetch cron job - runs every hour
0 * * * * cd /app && node scripts/fetch-repositories.js >> /var/log/cron/fetch-repositories.log 2>&1