# scripts/crontab
# Repository fetch cron job - runs daily at 2 AM
0 2 * * * cd /app && node scripts/fetch-repositories.js >> /var/log/cron/fetch-repositories.log 2>&1